# s6/envoy/run
#!/usr/bin/with-contenv sh
/usr/local/bin/envoy -c /etc/envoy/envoy.yaml